# accountindex/data/gearmand.conf
[program:gearmand]
command=/usr/sbin/gearmand -L 0.0.0.0 -l /var/log/gearmand.log')
numprocs=1
autostart=true
autorestart=true
user=root
umask=022
stopsignal=QUIT
redirect_stderr=true
