# accountindex/data/spamv2_front.conf
[program:spamv2front]
command=/usr/bin/python gearman_crawl.py -w

directory=/root/newsrc/accountindex/src/crawl/
numprocs=2
process_name=%(program_name)s_%(process_num)02d
umask=022
stopsignal=QUIT
stdout_logfile=/var/log/spamv2front_%(process_num)02d.log
redirect_stderr=true
